# catan-services_1.0/scripts/catan_webserver.conf
[program:catan_webserver]
command=/opt/catan/bin/catan_webserver
autostart=true
autorestart=true
stderr_logfile=/opt/catan/log/webserver.log.err
stdout_logfile=/opt/catan/log/webserver.log
